# ops/cron/marketing_warehouse.cron
# Template: replace <REPO_ROOT> with the absolute repository path.
# Example:
# 0 7 * * * /absolute/path/scripts/daily_refresh.sh >> /absolute/path/logs/cron_daily_refresh.log 2>&1
0 7 * * * <REPO_ROOT>/scripts/daily_refresh.sh >> <REPO_ROOT>/logs/cron_daily_refresh.log 2>&1
